# roles/sympa/templates/profile.d/sympa.sh.j2
# Required because mhonarc perl lib installed in non standard location
export PERL5LIB=/usr/share/mhonarc
